# build-docker-images/cron-backup/config/crontab
0 00 05 ? * * * /usr/bin/aws s3 sync /var/ossec/logs/alerts s3://S3_TMPL_PATH --only-show-errors
